# deploy/systemd/merge-stars-backend.service
[Unit]
Description=MERGE STARS NestJS API
After=network.target

[Service]
Type=simple
User=__DEPLOY_USER__
WorkingDirectory=__DEPLOY_DIR__/backend
EnvironmentFile=__DEPLOY_DIR__/.env
ExecStart=/usr/bin/node dist/main.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=merge-stars-backend

[Install]
WantedBy=multi-user.target
